# android_bridge_bot/rootfs/etc/s6-overlay/s6-rc.d/generate-config/up
#!/bin/sh
exec /etc/s6-overlay/scripts/generate-config.sh
