# android_bridge_bot/rootfs/etc/s6-overlay/s6-rc.d/abb/finish
#!/bin/sh
echo "[abb] exited with code $1, signal $2"
