# android_bridge_bot/rootfs/etc/s6-overlay/s6-rc.d/adb-server/finish
#!/bin/sh
echo "[adb-server] exited with code $1, signal $2"
/usr/bin/adb kill-server 2>/dev/null || true
